# src/macos/profile.d/linkedin.sh
#!/usr/bin/env bash

export PATH=/usr/local/linkedin/bin:$PATH
