# src/macos/profile.d/local.sh
#!/usr/bin/env bash

export PATH=/usr/local/bin:$(echo $PATH | sed 's/:\/usr\/local\/bin$//')
