# src/macos/profile.d/mysql.sh
#!/bin/bash
MYSQL_HOME=/usr/local/mysql/bin
PATH=$PATH:$MYSQL_HOME
export PATH
